# deployment/bitbot.service
[Unit]
Description=Bit-Bot Discord Service
# Wait for actual internet connection, not just networking subsystem
After=network-online.target
Wants=network-online.target

[Service]
# Replace 'youruser' with your actual Ubuntu username
User=youruser
Group=youruser
WorkingDirectory=/home/youruser/bitbot
ExecStart=/bin/bash /home/youruser/bitbot/deployment/start.sh

# Restart logic: Always restart, wait 10s between attempts
Restart=always
RestartSec=10

# Important: Don't stop trying to restart, even if it fails 50 times (long internet outage)
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
